# systemd.service
[Unit]
Description=MESERI Flask Application
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/meseri
Environment="PATH=/var/www/meseri/venv/bin"
Environment="FLASK_APP=app.py"
Environment="FLASK_ENV=production"
ExecStart=/var/www/meseri/venv/bin/gunicorn --config gunicorn.conf.py wsgi:app
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
